# etc/supervisord.d/dnsmasq.ini
[program:dnsmasq]
command=/opt/qnib/bin/start_dnsmasq.sh
autostart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log

